# luxbox_receiver.service
[Unit]
Description=Luxbox Receiver
After=network.target

[Service]
ExecStart=/opt/luxbox/luxbox_receiver.py

[Install]
WantedBy=multi-user.target
