# luxbox_sender.service
[Unit]
Description=Luxbox Sender
After=network.target

[Service]
ExecStart=/opt/luxbox/luxbox_sender.py

[Install]
WantedBy=multi-user.target
